# sbin/os/rc.local
#!/bin/sh -e
#
# rc.local for scatbot-pi

TARGET_DIR="/home/bee/scatbot"
LOG_FILE=$TARGET_DIR/logs/rclocal.log

cd $TARGET_DIR

echo "$(date) starting scatbot from rc.local" > $LOG_FILE
sudo bash ./start.sh >> $LOG_FILE 2>&1

exit 0
